# roles/proxmox_prepare/templates/zfs-textfile-exporter.service.j2
[Unit]
Description=Generate ZFS metrics for node_exporter textfile collector
After=zfs.target

[Service]
Type=oneshot
ExecStart={{ zfs_textfile_script_path }}